# rules/sigma/windows/registry/registry_run_key_persistence.yml
title: Registry Run Key Persistence
id: 1c4a6f7b-2a31-4212-b395-7f3e878b9410
status: experimental
description: Detects modification of common Windows Run and RunOnce registry keys used for persistence.
references:
  - https://attack.mitre.org/techniques/T1547/001/
author: Matt Caballero
date: 2026-07-09
modified: 2026-07-09
tags:
  - attack.persistence
  - attack.privilege_escalation
  - attack.t1547.001
logsource:
  product: windows
  category: registry_set
detection:
  selection_run_keys:
    TargetObject|contains:
      - '\\Software\\Microsoft\\Windows\\CurrentVersion\\Run'
      - '\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce'
  filter_known_paths:
    Details|contains:
      - 'Microsoft\\OneDrive'
      - 'Teams.exe'
  condition: selection_run_keys and not filter_known_paths
fields:
  - TargetObject
  - Details
  - Image
  - User
falsepositives:
  - Legitimate software installers and updaters adding startup entries.
level: medium
